# apps/kasko/.docker/supervisor.d/horizon.conf.default
[program:horizon]
process_name=%(program_name)s
command=php /var/www/kasko/artisan horizon
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/kasko/horizon.log
stopwaitsecs=3600
